# deploy/monolith/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
user=root

[program:gen_app_config]
command=/fba/gen_app_config.sh
autostart=true
autorestart=false
startsecs=0
priority=5
environment=VITE_DEMO_MODE="%(ENV_VITE_DEMO_MODE)s",DEMO_MODE="%(ENV_DEMO_MODE)s",VITE_APP_TITLE="%(ENV_VITE_APP_TITLE)s",VITE_GLOB_API_URL="%(ENV_VITE_GLOB_API_URL)s",VITE_TURNSTILE_SITE_KEY="%(ENV_VITE_TURNSTILE_SITE_KEY)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:fba_server]
directory=/fba/backend
command=/fba/.venv/bin/python -m granian backend.main:app --interface asgi --host 127.0.0.1 --port 8001 --workers 1 --backlog 1024 --workers-kill-timeout 120 --backpressure 2000 --log --log-level info
user=root
autostart=true
autorestart=true
startretries=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:fba_celery_worker]
directory=/fba/backend
# 使用 AsyncIOPool 代替 gevent，支持 async def 任务
# 必须设置 CELERY_CUSTOM_WORKER_POOL 环境变量
# environment=CELERY_CUSTOM_WORKER_POOL="celery_aio_pool.pool:AsyncIOPool"
# command=/fba/.venv/bin/python -m celery -A backend.app.task.celery worker --loglevel=INFO -P custom -c 4
# 使用 threads pool（兼容 asyncio.run()，适用于 async def 任务）
# gevent 不支持 asyncio，threads 可以在每个线程中运行 asyncio.run()
command=/fba/.venv/bin/python -m celery -A backend.app.task.celery worker --loglevel=INFO -P threads -c 4
user=root
autostart=true
autorestart=true
startretries=10
startsecs=10
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:fba_celery_beat]
directory=/fba/backend
command=/fba/.venv/bin/python -m celery -A backend.app.task.celery beat --loglevel=INFO
user=root
autostart=true
autorestart=true
startretries=10
priority=40
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
